# deploy/systemd/attendance-watchdog.timer
[Unit]
Description=Run Attendance System watchdog every minute
Requires=docker.service

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=attendance-watchdog.service

[Install]
WantedBy=timers.target
